# profile.d/alias.sh
alias ll='ls -al --color=auto'
alias pm='pnpm'